# clientmon/sfclientmon.conf
description "SolidFire client health monitor"
start on stopping dmesg
stop on runlevel [016]
respawn
chdir /opt/sfauto/clientmon
exec /usr/bin/python /opt/sfauto/clientmon/sfclientmon.py
