# services/sfhealthmon.conf
description "SolidFire cluster health monitor"
start on runlevel [2345]
stop on runlevel [016]
respawn
chdir /opt/lab/scripts/services
exec /usr/bin/python /opt/lab/scripts/services/sfhealthmon.py
