# services/vdbench.conf
description "Run vdbenchd"
# Start after firstboot has been run
start on stopping firstboot
stop on runlevel [016]
exec /opt/vdbench/vdbenchd.sh

