# client_daemons/vdbench_sfautovm.conf
description "Run vdbenchd after boot on a VM"
# Start after firstboot has been run
start on stopping firstboot
stop on runlevel [016]
exec /bin/bash /opt/sfauto/client_daemons/vdbenchd_sfautovm.sh
