# controller_daemons/sfclustermon.conf
description "SolidFire cluster monitor"
start on runlevel [2345]
stop on runlevel [016]
respawn
chdir /opt/sfclustermon
exec /usr/bin/python /opt/sfclustermon/sfclustermon.py

