# controller_daemons/sfevents.conf
description "Cluster events gathering"
start on stopping dmesg
stop on runlevel [016]
respawn
chdir /opt/sfevents/
exec /usr/bin/python /opt/sfevents/sfevents.py
