# controller_daemons/sfreporter.conf
description "Cluster reports gathering"
start on stopping dmesg
stop on runlevel [016]
respawn
chdir /opt/sfreporter
exec /usr/bin/python /opt/sfreporter/sfreporter.py
